# necessary.ini
[program:monitor]
command = /root/jm-nodelta-hft-system/.venv/bin/python /root/jm-nodelta-hft-system/radar/monitor.py
directory = /root/jm-nodelta-hft-system/radar
stdout_logfile = ./log/monitor.log
redirect_stderr = true
autorestart = true
stopsignal = TERM
stdout_logfile_maxbytes = 30MB
stdout_logfile_backups = 10
priority = 0